# robot_side/templates/agibot/m3v-agent-agibot.service
[Unit]
Description=Multi3DViz 受控端 agent (Agibot D1)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# ROS env so rclpy finds the distro. Source the setup.bash before exec.
Environment="PYTHONUNBUFFERED=1"
ExecStartPre=/bin/bash -c 'source /opt/ros/humble/setup.bash 2>/dev/null || source /opt/ros/noetic/setup.bash 2>/dev/null || true'
ExecStart=/bin/bash -c 'source /opt/ros/humble/setup.bash 2>/dev/null || source /opt/ros/noetic/setup.bash 2>/dev/null; exec python3 -m m3v_agent.agent --config /etc/m3v-agent/config.yaml'
WorkingDirectory=/opt/m3v-agent
Restart=on-failure
RestartSec=3
# Give the agent a chance to lie_down + flush on stop.
KillSignal=SIGTERM
TimeoutStopSec=10
# Run as root so we can write /home/<user>/... and bind UDP 43988 without friction.
# Harden in production: drop to a dedicated user + grant the needed caps.
User=root

[Install]
WantedBy=multi-user.target
